# etc/systemd/system/asusgpuboot.service
[Unit]
Description=ASUS G14 GPU boot script...
Before=graphical.target

[Service]
User=root
Type=simple
ExecStart=/usr/sbin/asus_gpu_boot

[Install]
WantedBy=graphical.target