# etc/systemd/system/asusgpuswitch.service
[Unit]
Description=ASUS G14 GPU switch script...

[Service]
User=root
Type=simple
ExecStart=/usr/bin/asusgpuswitch
