# etc/systemd/system/asus_backlight_off_before_suspend.service
[Unit]
Description=ASUS keyboard backlight off. . .
Before=suspend.target

[Service]
User=root
Type=simple
ExecStart=/usr/bin/asus_backlight_off

[Install]
WantedBy=suspend.target
