# etc/systemd/system/asusgpuboot.service
[Unit]
Description=dGPU boot enable/disable script...
Before=graphical.target
After=suspend.target

[Service]
User=root
Type=simple
ExecStart=/usr/bin/asusgpuboot

[Install]
WantedBy=graphical.target
WantedBy=suspend.target
